# openwrt-opennic.ssl/rc.local
pgrep uhttpd && /etc/init.d/uhttpd disable
pgrep uhttpd && /etc/init.d/uhttpd stop
